# etc/service/php/run
#!/bin/bash

cd /var/www/html

exec php-fpm -F 2>&1
